# provisioning/roles/supervisord/templates/conf.d/google-ime-skk.conf
[program:google-ime-skk]
command=/usr/local/bin/google-ime-skk -p 1178
autostart=true
autorestart=true
